# deploy/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/runpod-volume/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=3
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── Redis (embedded) ────────────────────────────────────────
[program:redis]
command=redis-server --port 6379 --save 60 1 --dir /runpod-volume/data --dbfilename redis-dump.rdb --loglevel warning
autostart=true
autorestart=true
priority=10
stdout_logfile=/runpod-volume/logs/redis.log
stderr_logfile=/runpod-volume/logs/redis-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

# ── Xvfb (virtual display for Chrome) ───────────────────────
[program:xvfb]
command=Xvfb :99 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
environment=DISPLAY=":99"
autostart=true
autorestart=true
priority=15
stdout_logfile=/runpod-volume/logs/xvfb.log
stderr_logfile=/runpod-volume/logs/xvfb-err.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

# ── Django Web Server (Gunicorn + Uvicorn on port 3000) ─────
[program:web]
command=gunicorn config.asgi:application
    -k uvicorn.workers.UvicornWorker
    --workers 2
    --bind 0.0.0.0:3000
    --graceful-timeout 30
    --timeout 0
    --access-logfile /runpod-volume/logs/gunicorn-access.log
    --error-logfile /runpod-volume/logs/gunicorn-error.log
environment=DISPLAY=":99"
directory=/app
autostart=true
autorestart=true
priority=30
startsecs=10
startretries=5
stdout_logfile=/runpod-volume/logs/web.log
stderr_logfile=/runpod-volume/logs/web-err.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB

# ── Celery Worker (browser tasks, event processing) ─────────
[program:celery-worker]
command=celery -A config worker -l info --pool=threads --concurrency=4
environment=DISPLAY=":99"
directory=/app
autostart=true
autorestart=true
priority=40
startsecs=10
startretries=5
stopwaitsecs=60
stdout_logfile=/runpod-volume/logs/celery-worker.log
stderr_logfile=/runpod-volume/logs/celery-worker-err.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB

# ── Celery Beat (scheduled tasks) ───────────────────────────
[program:celery-beat]
command=celery -A config beat --loglevel=info --scheduler=redbeat.RedBeatScheduler
directory=/app
autostart=true
autorestart=true
priority=50
startsecs=10
startretries=5
stdout_logfile=/runpod-volume/logs/celery-beat.log
stderr_logfile=/runpod-volume/logs/celery-beat-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
